# conf/supervisor/services.conf
[program:api-service]
command=python manage.py run
directory=/webapps
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/apps/api.log
